# src/deploy/supervisor/conf.d/tokumx.conf
[program:tokumx]
user=tokumx
command=/home/laisky/src/tokumx-1.5.0-linux-x86_64/bin/mongod --config=/etc/tokumx.conf
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/supervisor/tokumx.log
stderr_logfile=/var/log/supervisor/tokumx.err

